# plex-mount-recovery/files/check-media-mounts.service
[Unit]
Description=Check and mount media shares for Plex
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check-media-mounts.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
